# src/main/resources/etc/systemd/system/amhs.service
[Unit]

Description=AMHS

[Service]

ExecStart=/opt/amhs/bin/amhs

Restart=always

User=amhs

[Install]

WantedBy=multi-user.target